# fan-ctrl.service
[Unit]
Description=fanctrl

[Service]
Type=oneshot
ExecStart=/bin/sh /etc/init.d/fan-ctrl.sh

[Install]
WantedBy=multi-user.target
